# systemd/backup.timer
[Unit]
Description=Daily timer for the Secure Backup service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target