# systemd/backup.timer
[Unit]
Description=Temporizador per a la còpia de seguretat diària

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target